# docker/aws_tools/profile.d/ls-alias.sh
alias ls='ls -FA --color=tty'
alias l="ls"
alias ll="ls -l"
alias la="ls -a"
